# files/disable-tso
#!/bin/sh
### BEGIN INIT INFO
# Provides:          disable-tso
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: TSO disable script
# Description:       Disable TSO on Neutron interfaces
### END INIT INFO

# Authors: Emilien Macchi <emilien.macchi@enovance.com>
#          Sylvain Afchain <sylvain.afchain@enovance.com>

case "$1" in
start)
        if [ -f /var/run/disable-tso.pid ]; then
          "TSO script is already running or lock file has not been deleted."
          exit 1
        else
          ip link | awk -F ': ' '/qvb|qvo|qbr/{system("ethtool --offload "$2 " tso off gso off tx off ufo off gro off")}'
          ip monitor link | awk -F ': ' '/qvb|qvo|qbr/{system("ethtool --offload "$2 " tso off gso off tx off ufo off gro off")}'&
          touch /var/run/disable-tso.pid
          echo "Disabling TSO/GSO/GRO/UFO on QVB/QVO/QBR interfaces"
        fi
;;
stop)
        echo "Enabling TSO/GSO/GRO/UFO on QVB/QVO/QBR interfaces"
        if [ -f /var/run/disable-tso.pid ]; then
          pkill -xf '/sbin/ip monitor link'>&2
          rm /var/run/disable-tso.pid
        else
          "TSO script was not running or PID file is missing."
          exit 1
        fi
;;
*)
        echo "Usage: /etc/init.d/disable-tso {start|stop}"
        exit 1
;;
esac

exit 0
